# launchd/com.mikeybeez.llmos.kernel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mikeybeez.llmos.kernel</string>

    <!-- Named wrapper, not a bare python3 -> legible in Login Items / ps -->
    <key>ProgramArguments</key>
    <array>
        <string>/Users/bard/Code/LLMOS/bin/llmos-kernel</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <!-- Supervision comes free from launchd: restart the kernel if it dies -->
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/bard/Code/LLMOS/state/logs/kerneld.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/bard/Code/LLMOS/state/logs/kerneld.err.log</string>
</dict>
</plist>
